# build/sing-box-daemon.service
[Unit]
Description=sing-box reF1nd Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
StateDirectoryMode=0700
UMask=0077
StateDirectory=sing-box-daemon-reF1nd
ExecStart=/opt/sing-box-reF1nd/resources/daemon/sing-box-daemon run --working-directory /var/lib/sing-box-daemon-reF1nd --socket /run/sing-box-reF1nd.socket
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
